# contrib/kast.service
[Unit]
Description=Kast Service
Wants=display-manager.service kast-pre.service
After=network.target display-manager.service kast-pre.service

[Service]
ExecStart=/usr/bin/kast
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
